# dotfiles/easyeffects/easyeffects.service
[Unit]
Description=EasyEffects audio processing (background, Legion speaker EQ)
After=pipewire.service wireplumber.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/easyeffects --service-mode
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
