# camera/camera-privacy-guard.service
[Unit]
Description=Enforce the camera privacy switch by unbinding uvcvideo when privacy=ON
After=multi-user.target systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/camera-privacy-guard.sh
Restart=always
RestartSec=2
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/sys/bus/usb/drivers/uvcvideo

[Install]
WantedBy=multi-user.target
